# localhost/etc/supervisor.d/reverse tunnel.ini
[program:reverse-tunnel]
command=/bin/ssh-reverse-tunnel-daemon /home/daniel/reverse-tunnel-hosts
numprocs=1
autostart=true
startsecs=5
startretries=5
autorestart=true
stdout_logfile=/var/log/supervisor_reverse_tunnel.log
stderr_logfile=/var/log/supervisor_reverse_tunnel.error.log
